# pxe/s6/dnsmasq/run
#!/bin/bash
exec dnsmasq --no-daemon -C /etc/dnsmasq.d/pxe.conf
